# integration/triggerhappy/etc/systemd/system/triggerhappy.service.d/override.conf
[Service]
EnvironmentFile=/etc/triggerhappy.env
